# playbooks/roles/contrail/analytics/files/systemd/contrail-ansible.service
[Unit]
Description=Contrail Analytics configuration
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/contrailctl config sync -c analytics -F -t configure,service
TimeoutStartSec=0
MountFlags=slave

[Install]
WantedBy=multi-user.target
